# manager/upstart-zurl.conf
description "zurl"

start on runlevel [2345]
stop on runlevel [!2345]

limit nofile 20000 20000

respawn
exec /usr/local/bin/zurl

post-start script
	sleep 1
	chmod a+w /tmp/zurl-*
end script
